# automount/10_automount.rules.in
# Generated by @PACKAGE@ - DO NOT HAND EDIT!

#KERNEL=="sd*", RUN+="@sbindir@/automount"
KERNEL=="sd?[0-9]*", RUN+="@sbindir@/automount"
